# a9-cam/rootfs/etc/s6-overlay/s6-rc.d/init/finish
#!/command/with-contenv bashio
# shellcheck shell=bash

pkill -9 python3
