# a9-cam/rootfs/etc/s6-overlay/s6-rc.d/init/run
#!/command/with-contenv bashio
# shellcheck shell=bash

python3 /usr/local/a9-v720/src/a9_naxclow.py -s
